# other/omoide.service
[Unit]
Description=Omoide web server
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/Omoide
ExecStart=/home/pi/Omoide/run_server.sh
CPUAccounting=true
MemoryAccounting=true

[Install]
WantedBy=multi-user.target
